# dracut/96luks-2fa/luks-2fa.target
[Unit]
Description=2FA with LUKS
Before=sysinit.target dracut-pre-mount.service
After=systemd-ask-password-console.path
